# laravel/crontab/laravel.cron
* * * * * /usr/bin/php /apps/test/artisan schedule:run >> /dev/null 2>&1
